# health-monitor.service
[Unit]
Description=Simple Linux Disk Health Monitor
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/health-monitor/monitor.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target